# camera.path
[Unit]
Description=Monitor dev v4l for camera device

[Path]
PathExists=/dev/v4l/by-id/usb-046d_HD_Pro_Webcam_C920_5AB6901F-video-index0
Unit=split-camera.service

[Install]
WantedBy=multi-user.target